# Cron/cron.d/hooe-cron
* * * * * /usr/local/bin/php /usr/share/nginx/html/home/artisan schedule:run >> /usr/share/nginx/html/home/storage/logs/laravel.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job